# data/autostart/xfce4-panel.desktop
[Desktop Entry]
Type=Application
Exec=xfce4-panel
Name=xfce4-panel
OnlyShowIn=MATE;
