# data/autostart/xfwm4.desktop
[Desktop Entry]
Type=Application
Exec=xfwm4
Name=xfwm4
OnlyShowIn=MATE;
